# systemd/user/power-profile@ac.service
[Unit]
Description=Set performance profile when AC is plugged in

[Service]
Type=oneshot
ExecStart=/usr/bin/powerprofilesctl set performance
