# systemd/user/power-profile.service
[Unit]
Description=Automatically set power profile based on AC power

[Service]
Type=oneshot
ExecStart=/usr/bin/powerprofilesctl set %i
