# systemd/user/power-profile@battery.service
[Unit]
Description=Set power-saver profile when on battery

[Service]
Type=oneshot
ExecStart=/usr/bin/powerprofilesctl set power-saver

